# rancher-gce/nginx-ubuntu/certbot.cron
#get initial cert and then check for cert renewal everyhour
7 * * * *  /usr/local/certbot/certbot.sh > /tmp/certbot.sh.cron.log 2>&1;
