# rancher-gce/nginxssl-dockprox/cron/dockprox.cron
#create /etc/nginx/conf.d/ files and get certs or update if required
*/2 * * * * /var/local/certbot/certbot.sh;
