# rancher-gce/postfix-dockprox/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
priority=1
command=service cron start
autorestart=true
startsecs=60

[program:rsyslog]
priority=2
command=service rsyslog start
autorestart=true
startsecs=60

[program:postfix]
priority=10
command=service postfix start
autorestart=true
startsecs=60
